# nuxt-base/config/supervisor/conf.d/node-memory-monitor.conf
[program:node-memory-monitor]
command=sh -c "/scripts/node-memory-monitor.sh || kill -s SIGQUIT $PPID"

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

autorestart=false
startretries=0
stopasgroup=true
