# common/config/supervisor/conf.d/nginx.conf
[program:nginx]
command=sh -c "nginx -g 'daemon off;' || kill -s SIGQUIT $PPID"

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

autorestart=false
startretries=0
stopasgroup=true
